# curator4-1832765551.conf
[program:curator4-1832765551]
command=python3 /home/curator4/curator4-1832765551/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/curator4-1832765551/test.err.log
stdout_logfile=/home/curator4/curator4-1832765551/test.out.log

